# udev-rules/50-onbattery.rules
# Default to power-saver when on battery
SUBSYSTEM=="power_supply", ATTR{online}=="0", RUN+="/usr/local/bin/ppd-wrapper.sh set power-saver"
