# udev-rules/50-onbattery.rules
# Default to balanced when on battery
SUBSYSTEM=="power_supply", ATTR{online}=="0", RUN+="/usr/local/bin/ppd-wrapper.sh set balanced"
